# jhistd.service
[Unit]
Description=A jhist server daemon

[Service]
ExecStart=npm start
WorkingDirectory=/opt/jhistd
User=jhistd

[Install]
WantedBy=network.target
