# jhistd.service
[Unit]
Description=A jhist server daemon

[Service]
ExecStart=/usr/bin/node --harmony start
WorkingDirectory=/opt/jhistd
User=jhistd

[Install]
WantedBy=network.target
